# tv-playlist.service
[Unit]
Description=TV Playlist HTTP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/TV_BO
ExecStart=/usr/bin/python3 /root/TV_BO/serve.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
